# trading/deploy/trading.service
[Unit]
Description=hosub-trading (semi-auto trading dashboard)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=hosub
WorkingDirectory=/home/hosub/apps/hosub-mcp/trading
ExecStart=/home/hosub/apps/hosub-mcp/trading/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8600
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
